# deploy/supervisor/fundflow-queue.conf
; Laravel queue worker — background jobs (legacy import, notifications, etc.)
;
; Install on the app server:
;   sudo cp deploy/supervisor/fundflow-queue.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start fundflow-queue
;
; Requires .env: QUEUE_CONNECTION=database (or redis with matching connection below)

[program:fundflow-queue]
process_name=%(program_name)s
command=php /var/www/fundflow-saas/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --backoff=10
directory=/var/www/fundflow-saas
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/fundflow-saas/storage/logs/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600
