# deploy/supervisor/fundflow-reverb.conf
; Laravel Reverb — Filament / Echo WebSocket server
;
; Install on the app server:
;   sudo cp deploy/supervisor/fundflow-reverb.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start fundflow-reverb
;
; Requires .env: BROADCAST_CONNECTION=reverb, REVERB_SERVER_* (see docs/production-runbook.md)

[program:fundflow-reverb]
process_name=%(program_name)s
command=php /var/www/fundflow-saas/artisan reverb:start --no-interaction
directory=/var/www/fundflow-saas
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/fundflow-saas/storage/logs/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=15
